# backup/profile.d/cplus_path.sh
# Begin /etc/profile.d/cplus_path.sh

# Path for CPLUS_PATH if in the event error "could not find <stdlib.h>" occurs again
CPLUS_INCLUDE_PATH=/usr/include/c++/9.2.0:/usr/include

# End cplus_path.sh
