# backup/profile.d/mvn.sh
# Begin /etc/profile.d/mvn.sh

JAVA_HOME=/opt/jdk/bin

pathappend /opt/mvn/bin           PATH
#pathappend /opt/mvn/lib/pkgconfig PKG_CONFIG_PATH

export JAVA_HOME

# End /etc/profile.d/mvn.sh
